# supervisord-snmpd.conf
[program:snmpd]
command = /usr/sbin/snmpd
autorestart = true